# micromamba/etc/profile.d/mamba.sh.in
# Do not reset `MAMBA_ROOT_PREFIX` if it is already set
if [ -z "${MAMBA_ROOT_PREFIX}" ]; then
    export MAMBA_ROOT_PREFIX="@CMAKE_INSTALL_PREFIX@"
fi

__mamba_setup="$("@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/mamba" shell hook --shell posix 2> /dev/null)"
if [ $? -eq 0 ]; then
    eval "$__mamba_setup"
else
    alias mamba="@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/mamba"  # Fallback on help from mamba activate
fi
unset __mamba_setup
